# Files/rc.local
#!/bin/bash
export HOME=/root
export USER=root

date
ntpdate pool.ntp.org > /var/log/ntpdate.log 2>&1
crond

StartXPRA.sh
